# scripts/claude-prompt-viewer.service
[Unit]
Description=Claude Prompt Viewer - AI Coding Assistant Session Browser
After=network.target

[Service]
Type=simple
WorkingDirectory=%INSTALL_DIR%
ExecStart=/usr/bin/node %INSTALL_DIR%/node_modules/.bin/next start -p %PORT%
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=%PORT%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-prompt-viewer

[Install]
WantedBy=default.target
